# src/file/kde-baloo.service.in
[Unit]
Description=Baloo File Indexer Daemon
PartOf=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/baloo_file
BusName=org.kde.baloo
Slice=background.slice
ExecCondition=@KDE_INSTALL_FULL_BINDIR@/kde-systemd-start-condition --condition "baloofilerc:Basic Settings:Indexing-Enabled:true"

[Install]
WantedBy=graphical-session.target
